# setup/install/config_files/subscribers_fwd.conf
# subscribers - subscribers job file

description "Listens for uploaded surveys and sends results to subscribing systems"
author "Neil Penman <neilpenman@gmail.com>"

# Stanzas
#
# Stanzas control when and how a process is started and stopped
# See a list of stanzas here: http://upstart.ubuntu.com/wiki/Stanzas#respawn

# When to start the service
start on runlevel [2345]

# When to stop the service
stop on runlevel [016]

# Automatically restart process if crashed
#respawn

# Essentially lets upstart know the process will detach itself to the background
#expect fork

# Run before process
#pre-start script
#    [ -d /var/run/myservice ] || mkdir -p /var/run/myservice
#    echo "Put bash code here"
#end script

# Start the process
exec sudo -u tomcat8 /smap_bin/subscribers.sh default /smap forward
